# opensrp/composed/supervisord.conf
[supervisord]
nodaemon = true

[program:tomcat7]
command = /opt/tomcat/bin/catalina.sh run